# devices/z8803be/overlay/target/linux/mediatek/filogic/base-files/etc/uci-defaults/32-zbt-z8803be-wan-failover
#!/bin/sh
# SPDX-License-Identifier: GPL-2.0-only
#
# Z8803BE - dual wired-WAN topology + dormant WWAN DNS-policy contract.
#
# DNS POLICY CONTRACT (firmware-wide):
#   This board boots into deterministic public DNS by default. The
#   resolver path is:
#       LAN client -> dnsmasq (LAN IP:53) -> 1.1.1.1 / 8.8.8.8 (UDP/53)
#   This file owns the netifd-side half of that contract:
#       network.wan.peerdns      = 0   (drop ISP DHCP option 6)
#       network.wan_sfp.peerdns  = 0   (same policy for SFP WAN)
#       network.4_1.peerdns   = 0   (if cellular is explicitly enabled,
#                                    drop operator DNS pushed via netifd
#                                    for the qmodem proto=none stub)
#   Note: feeds/qmodem/.../modem_dial.sh:397 deletes
#   network.4_1.peerdns whenever the operator DNS list is empty
#   (Smart PH common case), so `uci show network` after a dial may
#   appear inconsistent with this comment. That's expected and not a
#   bug: noresolv=1 (below) makes the resolver path immune to
#   resolv.conf.auto state regardless.
#
#   The dnsmasq-side half lives in 80-zbt-z8803be-dns-cache, which
#   sets `noresolv=1` and pins `option server` to 1.1.1.1 / 8.8.8.8.
#   The qmodem-side half lives in 40-zbt-qmodem-dns-suppress, which
#   sets `do_not_add_dns=1` so quectel-CM-M never writes to
#   /etc/resolv.conf either.
#   Lastly, /etc/hotplug.d/iface/50-zbt-wwan-dns (the operator-DNS
#   capture hook) is shipped DISABLED via UCI (system.zbt_wwan_dns.
#   enabled=0) so it does not contradict the deterministic policy on
#   first boot. Setup scripts that prefer operator DNS can flip it on.
#
# WHY (root cause this fixes):
#   On a fresh `sysupgrade -n` first boot, the previous default
#   (peerdns=1) let udhcpc on the wired WAN interfaces write the ISP's
#   DNS into /tmp/resolv.conf.d/resolv.conf.auto BEFORE uci-defaults
#   could intervene. dnsmasq then forwarded queries to that ISP DNS in
#   round-robin alongside our seeded 1.1.1.1 / 8.8.8.8, and any
#   transient brokenness in operator DNS (Globe / Smart / Smart Bro PH
#   all observed) made the router boot with intermittent DNS even
#   though raw-IP routing worked. peerdns=0 here closes that hole at
#   the netifd level; noresolv=1 in the dnsmasq-side closes the
#   already-written-auto-file hole at the resolver level.

case "$(cat /tmp/sysinfo/board_name 2>/dev/null)" in
        zbtlink,zbt-z8803be|zbtlink,zbt-z8803be,mt7988a-nand) ;;
        *) exit 0 ;;
esac

set_slot_default_metric() {
	local target_slot="$1"
	local metric="$2"

	uci -q show qmodem 2>/dev/null | grep -E "^qmodem\.@modem-slot\[[0-9]+\]\.slot='${target_slot}'$" | \
		awk -F'.' '{print $2}' | awk -F'=' '{print $1}' | sort -u | while read -r sec; do
			[ -n "$sec" ] && uci -q set "qmodem.${sec}.default_metric=${metric}"
		done
}

uci -q set network.wan.metric='10'
uci -q set network.wan.defaultroute='1'
uci -q set network.wan.peerdns='0'
uci -q set network.wan.force_link='0'
if ! uci -q get network.wan6 >/dev/null 2>&1; then
	uci -q set network.wan6=interface
fi
uci -q set network.wan6.metric='10'
uci -q set network.wan6.device='@wan'
uci -q set network.wan6.proto='dhcpv6'

# Keep both wired uplinks DHCP-capable at first boot. Whichever link is
# actually present will come up automatically; if both are present, SFP
# wins by metric and the RJ45 WAN stays as the fallback path.
if uci -q get network.wan_sfp >/dev/null 2>&1; then
	uci -q set network.wan_sfp.proto='dhcp'
	uci -q set network.wan_sfp.metric='9'
	uci -q set network.wan_sfp.defaultroute='1'
	uci -q set network.wan_sfp.peerdns='0'
	uci -q set network.wan_sfp.force_link='0'
	if ! uci -q get network.wan_sfp6 >/dev/null 2>&1; then
		uci -q set network.wan_sfp6=interface
	fi
	uci -q set network.wan_sfp6.metric='9'
	uci -q set network.wan_sfp6.device='@wan_sfp'
	uci -q set network.wan_sfp6.proto='dhcpv6'
fi

# Register dormant WWAN (4_1 / 4_1v6) proto=none stubs bound to wwan0.
# Firmware defaults to wired WAN only: cellular has high metric, no
# default route, and IPv6 disabled. If a deployment explicitly enables
# cellular later, these stubs still let fw4 bind wwan0 to the wan zone.
if ! uci -q get network.4_1 >/dev/null 2>&1; then
	uci -q set network.4_1=interface
	uci -q set network.4_1.proto='none'
	uci -q set network.4_1.device='wwan0'
fi
uci -q set network.4_1.metric='200'
uci -q set network.4_1.defaultroute='0'
uci -q set network.4_1.peerdns='0'
# Mirror device → ifname (legacy alias). qmodem_monitor's
# update_netcfg() reads `network.<sec>.ifname` (modem_monitor.sh:81)
# but modern netifd canonicalises on `device`. Without this, the
# curl/ping probes run with --interface "" and the carrier-side
# liveness check silently fails. Setting both keeps qmodem_monitor
# happy without breaking netifd (it accepts either).
uci -q set network.4_1.ifname='wwan0'

if ! uci -q get network.4_1v6 >/dev/null 2>&1; then
	uci -q set network.4_1v6=interface
	uci -q set network.4_1v6.proto='none'
	uci -q set network.4_1v6.device='wwan0'
fi
uci -q set network.4_1v6.metric='200'
uci -q set network.4_1v6.disabled='1'

if [ -e /etc/config/qmodem ]; then
	set_slot_default_metric '4-1' '200'
	if uci -q get qmodem.4_1 >/dev/null 2>&1; then
		uci -q set qmodem.4_1.metric='200'
	fi
	uci -q commit qmodem
fi

uci -q commit network

# Keep dormant WWAN stubs and SFP WAN listed in the wan firewall zone.
# qmodem's boot logic has been observed wiping firewall.<wan zone>.network
# back to just 'wan' (it appears to rewrite the file when it (re)creates
# its own network sections), which silently drops wwan0 out of the wan
# zone if cellular is later enabled. add_list is idempotent because UCI
# deduplicates, so this is safe to re-run.
#
# Resolve the wan zone by NAME instead of by @zone[1] index: a user who
# adds a guest zone via LuCI before wan would shift the index, and we'd
# silently mutate the wrong zone.
find_wan_zone() {
	local sec
	for sec in $(uci -q show firewall 2>/dev/null \
		| sed -nE 's/^firewall\.(@zone\[[0-9]+\])=zone$/\1/p'); do
		if [ "$(uci -q get "firewall.${sec}.name")" = "wan" ]; then
			echo "$sec"
			return 0
		fi
	done
	return 1
}

if wan_zone="$(find_wan_zone)"; then
        zone_networks=$(uci -q get "firewall.${wan_zone}.network")
        case " $zone_networks " in
                *" wan_sfp "*) ;;
                *) uci -q add_list "firewall.${wan_zone}.network=wan_sfp" ;;
        esac
        case " $zone_networks " in
                *" wan_sfp6 "*) ;;
                *) uci -q add_list "firewall.${wan_zone}.network=wan_sfp6" ;;
        esac
        case " $zone_networks " in
                *" 4_1 "*) ;;
                *) uci -q add_list "firewall.${wan_zone}.network=4_1" ;;
        esac
        case " $zone_networks " in
                *" 4_1v6 "*) ;;
                *) uci -q add_list "firewall.${wan_zone}.network=4_1v6" ;;
        esac
        uci -q commit firewall
fi

# DNS-policy: seed the operator-DNS-capture hotplug as DISABLED so it
# does not write a `server=...` snippet into /tmp/dnsmasq.*.d/ on
# WWAN ifup. With noresolv=1 and pinned 1.1.1.1 / 8.8.8.8 forwarders
# (see 80-zbt-z8803be-dns-cache), letting that hotplug fire on first
# boot would re-introduce the exact race condition this commit set
# fixes. Setup scripts that prefer operator DNS can flip
# system.zbt_wwan_dns.enabled to '1' AFTER understanding the trade-off.
#
# We use a named uci section (zbt_wwan_dns) instead of an anonymous
# one so the toggle is greppable from setup scripts and LuCI without
# having to chase @section[N] index drift. /etc/config/system is
# always present (procd creates it), so we skip the existence guard.
if ! uci -q get system.zbt_wwan_dns >/dev/null 2>&1; then
        uci -q set system.zbt_wwan_dns=zbt_wwan_dns
        uci -q set system.zbt_wwan_dns.enabled='0'
        uci -q commit system
fi

exit 0
