# devices/z8803be/overlay/target/linux/mediatek/filogic/base-files/etc/uci-defaults/34-zbt-z8803be-usb-tether
#!/bin/sh

case "$(cat /tmp/sysinfo/board_name 2>/dev/null)" in
	zbtlink,zbt-z8803be|zbtlink,zbt-z8803be,mt7988a-nand) ;;
	*) exit 0 ;;
esac

uci -q set network.usb_tether=interface
uci -q set network.usb_tether.proto='dhcp'
uci -q set network.usb_tether.metric='50'
uci -q set network.usb_tether.defaultroute='1'
uci -q set network.usb_tether.peerdns='0'
uci -q commit network

find_wan_zone() {
	local sec
	for sec in $(uci -q show firewall 2>/dev/null | sed -nE 's/^firewall\.(@zone\[[0-9]+\])=zone$/\1/p'); do
		if [ "$(uci -q get "firewall.${sec}.name")" = "wan" ]; then
			echo "$sec"
			return 0
		fi
	done
	return 1
}

if wan_zone="$(find_wan_zone)"; then
	zone_networks=$(uci -q get "firewall.${wan_zone}.network")
	case " $zone_networks " in
		*" usb_tether "*) ;;
		*) uci -q add_list "firewall.${wan_zone}.network=usb_tether" ;;
	esac
	uci -q commit firewall
fi

if [ -x /etc/init.d/usbmuxd ]; then
	/etc/init.d/usbmuxd enable >/dev/null 2>&1 || true
	/etc/init.d/usbmuxd start >/dev/null 2>&1 || true
fi

exit 0
