# devices/z8803be/overlay/target/linux/mediatek/filogic/base-files/etc/uci-defaults/76-zbt-z8803be-admin-password
#!/bin/sh
# SPDX-License-Identifier: GPL-2.0-only
#
# Seed the root account with the vendor-label default password "admin".
# This single account is used by both LuCI (web UI) and dropbear (SSH).
#
# Runs only on Z8803BE, only once (script is deleted from
# /etc/uci-defaults by procd after a clean exit), and only when the
# root password is still empty - so subsequent manual changes via
# `passwd` or LuCI are preserved across later sysupgrades that happen
# to include this same script.

case "$(cat /tmp/sysinfo/board_name 2>/dev/null)" in
	zbtlink,zbt-z8803be|zbtlink,zbt-z8803be,mt7988a-nand) ;;
	*) exit 0 ;;
esac

# Probe the current root password hash. OpenWrt's /etc/shadow format:
#   root:<hash>:...
# Empty hash (field = "" or "!") means "no password set yet".
current_hash=$(awk -F: '$1 == "root" { print $2 }' /etc/shadow 2>/dev/null)
case "$current_hash" in
	""|"!"|"*")
		# No password set - apply the vendor default.
		printf 'admin\nadmin\n' | passwd root >/dev/null 2>&1
		;;
	*)
		# User already set one - leave it alone.
		;;
esac

exit 0
