# devices/z8803be/overlay/target/linux/mediatek/filogic/base-files/etc/uci-defaults/88-zbt-wrtbwmon-netdev-migrate
#!/bin/sh
# One-shot migration: wrtbwmon moved from `inet fw4` (forward-chain
# dispatch, undercount on offloaded flows) to `netdev wrtbwmon_acct`
# (ingress + egress hooks at priority -300 on every br-lan member).
#
# This script:
#   1. Rewrites /etc/wrtbwmon.conf so NFT_TABLE points at the new
#      netdev table. Sysupgrade preserves /etc/wrtbwmon.conf verbatim,
#      so on first boot after an upgrade we still see the old
#      `inet fw4` value and chain_manager.sh writes egress hooks into
#      the legacy table — which fails because inet does not support
#      egress hooks. Rewriting up front avoids the noisy first-run
#      failure.
#   2. Wipes any wrtbwmon-* leftovers from `inet fw4` so a router
#      upgraded from an older firmware doesn't keep stale chains, maps,
#      or per-domain sets alongside the new netdev table.
# The new table itself is created by `wrtbwmon-control apply` and
# `wrtbwmon setup`, so we don't pre-create anything here.

# --- 1. config rewrite (preserved /etc/wrtbwmon.conf may still point at fw4) ---
WCONF=/etc/wrtbwmon.conf
if [ -f "$WCONF" ] && grep -q '^NFT_TABLE="inet fw4"' "$WCONF"; then
    sed -i 's|^NFT_TABLE="inet fw4"|NFT_TABLE="netdev wrtbwmon_acct"|' "$WCONF"
    logger -t zbt-wrtbwmon-migrate "Rewrote NFT_TABLE in $WCONF to netdev wrtbwmon_acct"
fi

LEGACY="inet fw4"

# fw4 must already be loaded (uci-defaults runs after firewall start)
nft list table $LEGACY >/dev/null 2>&1 || exit 0

# Drop any wrtbwmon-dispatch rules from the forward chain.
nft list chain $LEGACY forward -a 2>/dev/null | \
    awk '/wrtbwmon-dispatch/ { print $NF }' | \
    while read -r handle; do
        nft delete rule $LEGACY forward handle "$handle" 2>/dev/null || true
    done

# Drop any leftover device chains.
nft list table $LEGACY 2>/dev/null | \
    awk '/^[[:space:]]*chain (device_|device_domains_)/ { gsub("{",""); print $2 }' | \
    while read -r ch; do
        nft flush  chain $LEGACY "$ch" 2>/dev/null || true
        nft delete chain $LEGACY "$ch" 2>/dev/null || true
    done

# Drop dispatch maps.
nft delete map $LEGACY wrtbwmon_dispatch_v4 2>/dev/null || true
nft delete map $LEGACY wrtbwmon_dispatch_v6 2>/dev/null || true

# Drop per-domain sets (named d_<MAC>_<DOMAIN>[_dl]).
nft list table $LEGACY 2>/dev/null | \
    awk '/^[[:space:]]*set d_[0-9a-f]+_/ { gsub("{",""); print $2 }' | \
    while read -r s; do
        nft delete set $LEGACY "$s" 2>/dev/null || true
    done

logger -t zbt-wrtbwmon-migrate "Migrated wrtbwmon from inet fw4 to netdev wrtbwmon_acct"
exit 0
