# devices/z8803be/overlay/package/luci-app-wifi-clients/root/etc/uci-defaults/90-zbt-wifi-clients-enable
#!/bin/sh

[ -x /etc/init.d/wifi_clients ] || exit 0
/etc/init.d/wifi_clients enable >/dev/null 2>&1 || true
/etc/init.d/wifi_clients start >/dev/null 2>&1 || true

exit 0
