# devices/z8803be/overlay/package/luci-app-zbt-modem-events/root/etc/init.d/zbt_modem_events
#!/bin/sh /etc/rc.common

START=99
STOP=10

CRON_FILE=/etc/crontabs/root
CRON_CMD='/usr/sbin/zbt-modem-events sample >/dev/null 2>&1'

remove_cron_line() {
	local tmp
	mkdir -p /etc/crontabs
	touch "$CRON_FILE"
	tmp=$(mktemp /tmp/zbt-modem-events-cron.XXXXXX)
	grep -Fv "$CRON_CMD" "$CRON_FILE" > "$tmp" || true
	cat "$tmp" > "$CRON_FILE"
	rm -f "$tmp"
}

add_cron_line() {
	remove_cron_line
	printf '*/5 * * * * %s\n' "$CRON_CMD" >> "$CRON_FILE"
}

reload_cron() {
	/etc/init.d/cron restart >/dev/null 2>&1 || true
}

start() {
	add_cron_line
	/usr/sbin/zbt-modem-events sample >/dev/null 2>&1 || true
	reload_cron
}

stop() {
	remove_cron_line
	reload_cron
}

restart() {
	stop
	start
}

reload() {
	start
}
